# cookbooks/opengenera/files/default/rc.local
#!/bin/sh

# Clear VNC pidfile on startup
rm -f /root/.vnc/genera-host:1.pid